# scripts/systemd/roxx-web.service
[Unit]
Description=RoXX Web Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=roxx
Group=roxx
WorkingDirectory=/opt/roxx
Environment=ROXX_CONFIG_DIR=/etc/roxx
Environment=ROXX_DATA_DIR=/var/lib/roxx
Environment=ROXX_LOG_DIR=/var/log/roxx
ExecStart=/opt/roxx/roxx server
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
